# system_files/shared/usr/lib/systemd/user/blossomos-flatpak-overrides-user.service
[Unit]
Description=BlossomOS Flatpak Overrides (per-user installs)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/blossomos-flatpak-overrides --user

[Install]
WantedBy=default.target
