# system_files/shared/usr/lib/systemd/system/blossomos-flatpak-overrides.service
[Unit]
Description=BlossomOS Flatpak Overrides
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/blossomos-flatpak-overrides --system

[Install]
WantedBy=multi-user.target
