# system_files/shared/usr/lib/systemd/system/blossomos-fanctrl-pause.service
[Unit]
Description=BlossomOS Pause fw-fanctrl during USB-C/USB4 partner negotiation
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/libexec/blossomos-fanctrl-pause
